# deploy/systemd/majlis.service
[Unit]
Description=Majlis — Multi-LLM Discussion Arena
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/majlis/repo

EnvironmentFile=/opt/majlis/repo/.env

ExecStartPre=/usr/bin/docker compose pull --quiet backend nginx
ExecStart=/usr/bin/docker compose up -d
ExecStop=/usr/bin/docker compose down

# Give vLLM models time to load on startup
TimeoutStartSec=600
TimeoutStopSec=60

Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
